# deploy/postavka-bot.service
[Unit]
Description=Postavka Assistant Bot (Telegram, FBO Ozon)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=postavka
Group=postavka
WorkingDirectory=/opt/postavka-bot
EnvironmentFile=/opt/postavka-bot/.env
ExecStart=/opt/postavka-bot/.venv/bin/python -m src.bot.main
Restart=always
RestartSec=10
TimeoutStopSec=20
KillSignal=SIGTERM

# Hardening (изоляция через отдельного юзера postavka — основная защита).
# ProtectSystem=strict не используется: pathlib.mkdir(exist_ok=True) всё равно
# делает syscall и ловит EROFS даже на существующей папке (logs/).
NoNewPrivileges=true
PrivateTmp=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true

# Logs -> journald
StandardOutput=journal
StandardError=journal
SyslogIdentifier=postavka-bot

[Install]
WantedBy=multi-user.target
